# config/flag-display.service
[Unit]
Description=E-Ink Flag Display Updater
After=network-online.target
Wants=network-online.target

[Service]
User=chris
WorkingDirectory=/home/chris/Flags
ExecStart=/usr/bin/python3 /home/chris/Flags/scripts/main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target